# scripts/medcat.service
[Unit]
Description=MedCATtrainer service

[Service]
ExecStart=/home/ubuntu/scripts/start_med_cat_trainer.sh

[Install]
WantedBy=multi-user.target

